# deploy/systemd/lumen-hub.service
[Unit]
# Lumen hub — single-binary monitoring server.
# See: docs/install/hub-binary.md for a walkthrough.
Description=Lumen hub (monitoring server)
Documentation=https://github.com/quanla93/lumen
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lumen
Group=lumen
EnvironmentFile=/etc/lumen/hub.env
WorkingDirectory=/var/lib/lumen
ExecStart=/usr/local/bin/lumen-hub

# Restart on crash (transient) or non-zero exit (config error). Don't
# restart on a clean exit so `systemctl stop` actually stops.
Restart=on-failure
RestartSec=5s

# Hardening — the hub only needs to read /etc/lumen/hub.env, write to
# /var/lib/lumen (SQLite), and listen on a port. Everything else gets
# locked down so a hub compromise is contained.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
ReadWritePaths=/var/lib/lumen
# Allow the hub to listen on its configured port. If you change LUMEN_HUB_ADDR
# to a port < 1024 you'll need to add CAP_NET_BIND_SERVICE or use a reverse
# proxy (recommended).
CapabilityBoundingSet=
AmbientCapabilities=

# Resource limits — hub stays well under 100 MB even with 100 hosts.
MemoryMax=512M
TasksMax=128

[Install]
WantedBy=multi-user.target
